# .agent/skills/context-engineering-collection/researcher/orchestration/launchd/com.context-engineering.loop-step.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.context-engineering.loop-step</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/researcher/orchestration/launchd/run-loop-step.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/researcher/reports/logs/launchd-loop-step.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/researcher/reports/logs/launchd-loop-step.err</string>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
