# .agent/skills/context-engineering-collection/researcher/orchestration/launchd/com.context-engineering.loop-discover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.context-engineering.loop-discover</string>
    <key>ProgramArguments</key>
    <array>
        <string>__REPO_ROOT__/researcher/orchestration/launchd/run-loop-discover.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key>
            <integer>5</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Hour</key>
            <integer>17</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__REPO_ROOT__/researcher/reports/logs/launchd-loop-discover.out</string>
    <key>StandardErrorPath</key>
    <string>__REPO_ROOT__/researcher/reports/logs/launchd-loop-discover.err</string>
    <key>WorkingDirectory</key>
    <string>__REPO_ROOT__</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
